# deploy/outreach-api.service
[Unit]
Description=Outreach Action API
After=network.target

[Service]
User=ubuntu
WorkingDirectory=/home/ubuntu/Scraper
Environment=PATH=/home/ubuntu/.local/bin:/usr/bin:/bin
ExecStart=/home/ubuntu/.local/bin/uv --project backend run python backend/action_api.py
Restart=always
RestartSec=5

[Install]
WantedBy=multi-user.target
