# deploy/outreach-listener.service
[Unit]
Description=Outreach Telegram Listener
After=network.target

[Service]
User=ubuntu
WorkingDirectory=/home/ubuntu/Scraper
Environment=PATH=/home/ubuntu/.local/bin:/usr/bin:/bin
ExecStart=/home/ubuntu/.local/bin/uv --project backend run python backend/telegram_bot.py
Restart=always
RestartSec=5

[Install]
WantedBy=multi-user.target
